# scheduler/launchd/com.hawkstrade.crypto-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hawkstrade.crypto-scan</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/arunbabuchandrababu/Desktop/AIPROJECTS/HawksTrade/scheduler/launchd/hawkstrade_launchd_runner.sh</string>
    <string>crypto-scan</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/arunbabuchandrababu/Desktop/AIPROJECTS/HawksTrade</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/arunbabuchandrababu/Desktop/AIPROJECTS/HawksTrade/logs/launchd_crypto_scan.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/arunbabuchandrababu/Desktop/AIPROJECTS/HawksTrade/logs/launchd_crypto_scan.err</string>
</dict>
</plist>
