# scheduler/systemd/hawkstrade-health-check.timer
[Unit]
Description=Run HawksTrade health checks every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
AccuracySec=1min
Unit=hawkstrade-health-check.service

[Install]
WantedBy=timers.target
