# scheduler/systemd/hawkstrade-crypto-scan.timer
[Unit]
Description=Run HawksTrade crypto-only scans hourly

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=1min
Unit=hawkstrade-crypto-scan.service

[Install]
WantedBy=timers.target
