# scheduler/systemd/hawkstrade-risk-check.timer
[Unit]
Description=Run HawksTrade risk checks during weekday market hours

[Timer]
OnCalendar=Mon..Fri *-*-* 13:45:00
OnCalendar=Mon..Fri *-*-* 14..19:00:00
OnCalendar=Mon..Fri *-*-* 14..19:15:00
OnCalendar=Mon..Fri *-*-* 14..19:30:00
OnCalendar=Mon..Fri *-*-* 14..19:45:00
Persistent=true
AccuracySec=1min
Unit=hawkstrade-risk-check.service

[Install]
WantedBy=timers.target
