# scheduler/systemd/hawkstrade-secrets.service
[Unit]
Description=HawksTrade /dev/shm secret loader
Documentation=file:/home/ec2-user/HawksTrade/scheduler/systemd/README.md
Wants=network-online.target
After=network-online.target
Before=hawkstrade-stock-scan.service hawkstrade-full-scan.service hawkstrade-crypto-scan.service hawkstrade-risk-check.service hawkstrade-daily-report.service hawkstrade-weekly-report.service hawkstrade-health-check.service
ConditionPathExists=/etc/hawkstrade/hawkstrade.secrets

[Service]
Type=oneshot
EnvironmentFile=/etc/hawkstrade/hawkstrade.env
ExecStart=/usr/bin/bash -lc 'install -m 0640 -o root -g "${HAWKSTRADE_GROUP:-ec2-user}" "${HAWKSTRADE_SECRET_SOURCE:-/etc/hawkstrade/hawkstrade.secrets}" "${HAWKSTRADE_SHM_SECRET_FILE:-/dev/shm/.hawkstrade.env}" && test -s "${HAWKSTRADE_SHM_SECRET_FILE:-/dev/shm/.hawkstrade.env}"'

[Install]
WantedBy=multi-user.target
